# host-assets/systemd/ltfs-lto6.service
[Unit]
Description=LTFS mount for HP LTO-6
After=local-fs.target network.target
Requires=local-fs.target
Conflicts=tape-safe-eject.service ltfs-deep-recovery.service

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/ltfs-lto6
ExecStart=/usr/bin/python3 /usr/local/tools/ltfs_recovery.py --orchestrated-mount
ExecStop=/usr/bin/python3 /usr/local/tools/ltfs_recovery.py --orchestrated-stop
TimeoutStartSec=480
TimeoutStopSec=120
KillMode=mixed
KillSignal=SIGTERM
FinalKillSignal=SIGKILL

[Install]
WantedBy=multi-user.target